# jupyter.timer
[Unit]
Description=Jupyter timer

[Timer]
OnUnitActiveSec=120s
OnBootSec=120s

[Install]
WantedBy=timers.target